# systemd/llm-interface.service
[Unit]
Description=Local LLM Chat Interface
After=default.target

[Service]
Type=simple
WorkingDirectory=@WORKING_DIRECTORY@
Environment=PORT=17865
Environment=OLLAMA_BASE_URL=http://127.0.0.1:11434
Environment=OLLAMA_MODEL=mistral
ExecStart=/usr/bin/env @EXEC_START@
Restart=on-failure
RestartSec=3
NoNewPrivileges=yes
PrivateDevices=yes
PrivateTmp=yes
ProtectControlGroups=yes
ProtectHome=read-only
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectSystem=strict
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
CapabilityBoundingSet=
AmbientCapabilities=
SystemCallArchitectures=native
UMask=0077

[Install]
WantedBy=default.target
